# config/supervisor/goserver.ini
[program:server]
command=/var/app/server
redirect_stderr=true