# config/supervisor/haproxy.ini
[program:server]
command=/usr/sbin/haproxy -f /etc/haproxy/haproxy.cfg
redirect_stderr=true